# supervisor/spider-server.ini
[program:cunxun-poster]
command=/root/www/work/gopath/bin/phantomjs spider-server.js 8006
directory=/data/apps/html/wangqingang/cunxun-poster
autostart=true
autorestart=true
startsecs=2
startretries=3
stdout_logfile=/data/logs/supervisor/cunxun-poster/access.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=20
stderr_logfile=/data/logs/supervisor/cunxun-poster/stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=2
